# post_install_overlay/local.d/alive.start
#! /bin/bash
echo It's alive >> /var/tmp/alive.log
